# templates/upstart/foreman/master.conf.erb
start on vagrant-mounted or runlevel [2345]
stop on runlevel [!2345]
